# ironman_v1/provisioners/files/nginx-supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/bin/bash -c "service nginx restart"
